# scripts/clawos-watchdog.timer.template
[Unit]
Description=Run ClawOS Watchdog every 10 minutes

[Timer]
OnBootSec=5min
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target
